# media-kit-api.service
[Unit]
Description=Media Kit Search API
After=network.target

[Service]
Type=exec
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/deep-brochures
Environment=PATH=/home/ubuntu/deep-brochures/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=VIRTUAL_ENV=/home/ubuntu/deep-brochures/venv
ExecStart=/home/ubuntu/deep-brochures/venv/bin/python -m gunicorn main:app -c gunicorn.conf.py
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 